# piops-fleet.service
[Unit]
Description=PiOps ESP32 Fleet Dashboard
After=network.target mosquitto.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/piops_dashboard
ExecStart=/usr/bin/python3 /home/pi/piops_dashboard/app.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
